# packaging/common/systemd/cliprelay-wl-watch.service
[Unit]
Description=ClipRelay Wayland clipboard watch

[Service]
Type=simple
EnvironmentFile=-%h/.config/cliprelay/cliprelay.env

# Defaults (can be overridden by EnvironmentFile)
Environment=CLIPRELAY_RELAY=127.0.0.1:8080
Environment=CLIPRELAY_ROOM=default
Environment=CLIPRELAY_MAX_TEXT_BYTES=1048576
Environment=CLIPRELAY_MAX_IMAGE_BYTES=20971520
Environment=CLIPRELAY_MAX_FILE_BYTES=20971520
Environment=CLIPRELAY_IMAGE_MODE=force-png
Environment=CLIPRELAY_WATCH_MODE=watch
Environment=CLIPRELAY_POLL_INTERVAL_MS=200

ExecStart=cliprelay-node wl-watch \
	--room ${CLIPRELAY_ROOM} \
	--relay ${CLIPRELAY_RELAY} \
	--mode ${CLIPRELAY_WATCH_MODE} \
	--interval-ms ${CLIPRELAY_POLL_INTERVAL_MS} \
	--max-text-bytes ${CLIPRELAY_MAX_TEXT_BYTES} \
	--max-image-bytes ${CLIPRELAY_MAX_IMAGE_BYTES} \
	--max-file-bytes ${CLIPRELAY_MAX_FILE_BYTES} \
	--image-mode ${CLIPRELAY_IMAGE_MODE}
Restart=on-failure
RestartSec=1

[Install]
WantedBy=default.target
